# netdata/rootfs/etc/services.d/netdata/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: netdata
# Runs netdata script
# ==============================================================================


bashio::log.info "Starting netdata..."

exec /usr/sbin/netdata
